# daily_updater/crontab
0 22 * * * python /code/src/daily_updater/main.py 